# cron/update-force-reboot.cron
# update the system and force a reboot after the update
@daily root update --auto-clean-log --reboot
